# uptime-kuma-monitoring/monitor-internet.timer.j2
[Unit]
Description=Internet monitoring timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
