# uptime-kuma-monitoring/monitor-internet.service.j2
[Unit]
Description=Internet Monitor and report status to Uptime Kuma

[Service]
Type=oneshot
ExecStart={{ script_path }}
